# services/bitlbee/run
[ ! -d "/run/bitlbee/" ] && mkdir -p /run/bitlbee
exec bitlbee -D
